# ncsc_feed.timer
[Unit]
Description=NCSC Feed updater
After=syslog.target

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target

